# scripts/web-app.service
[Unit]
Description=go-web-app

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/bin/web-app

[Install]
WantedBy=multi-user.target